# inc/supervisor.conf
[program:willemallan]
command=/home/willemallan/env/bin/python /home/willemallan/project/manage.py run_gunicorn -b 127.0.0.1:8000 -w 3
directory=/home/willemallan/project
autostart=true
autorestart=true